# folio.conf
[program:folio]
directory=/home/ubuntu/bots/folio
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/bots/folio/logfile_err.log
stdout_logfile=/home/ubuntu/bots/folio/logfile.log



